# systemd/pi-fpv-companion.service
[Unit]
Description=pi-fpv-companion drone visual tracker
After=multi-user.target network.target
# Recovery ladder (StartLimit* MUST live in [Unit] — ignored in [Service]).
# Rung 1: the frame-watchdog os._exit()s on a stall (~2s) / no-first-frame (~15s).
# Rung 2: systemd restarts the process (RestartSec=1) — clears a libcamera/IMX500
#         capture hang in ~3-5s by reopening the camera (bench: 10/10, no hang).
# Rung 3: ExecStartPre (camera-ready.sh) — if the sensor fell off the bus (not just a
#         capture hang), reset the camera subsystem via driver unbind/rebind, NO reboot.
# Rung 4: StartLimitAction=reboot if rung 3 keeps failing — full Pi reboot, last resort.
#         A reboot fully re-inits the CSI/sensor. The companion runs STANDBY-only and
#         the pilot keeps the FPV cam, so a ~48s outage is tolerable. =none to disable.
StartLimitIntervalSec=180
StartLimitBurst=8
StartLimitAction=reboot

[Service]
Type=simple
User=pi
Group=video
SupplementaryGroups=dialout

WorkingDirectory=/opt/pi-fpv-companion
Environment=PYTHONUNBUFFERED=1
# Wait for the IMX500 to be enumerable before launching; if it isn't, reset the camera
# subsystem (driver rebind, no reboot) — recovery rung 3. The '+' runs it as root
# (sysfs unbind/rebind needs root) regardless of User=. See scripts/camera-ready.sh.
ExecStartPre=+/opt/pi-fpv-companion/scripts/camera-ready.sh
ExecStart=/opt/pi-fpv-companion/.venv/bin/python -m pi_fpv_companion \
    --config /opt/pi-fpv-companion/config/imx500.yaml

Restart=on-failure
RestartSec=1
TimeoutStopSec=10

# Filesystem hardening — drone software shouldn't be writing outside its install
ProtectSystem=strict
ReadWritePaths=/opt/pi-fpv-companion/var
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

# Resource caps — fail loud if the pipeline tries to balloon
MemoryMax=300M
TasksMax=64

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-fpv-companion

[Install]
WantedBy=multi-user.target
